# infra/launchd/com.asif.babu-journal-static.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Auto-starts the Babu Journal static file server (serves site/) at user login.

  BEFORE INSTALLING:
    1. Replace {{REPO}} below with the absolute path to the Journal repo,
       e.g. /Users/asifhussain/PROJECTS/journal
    2. Copy the result to ~/Library/LaunchAgents/
    3. Run: launchctl load -w ~/Library/LaunchAgents/com.asif.babu-journal-static.plist

  The -l flag on bash ensures the login profile is read (picks up nvm, Homebrew node paths).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.asif.babu-journal-static</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd {{REPO}}/site &amp;&amp; exec npx serve . -l 3000 --cors</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>{{REPO}}/server/.logs/static.out.log</string>

  <key>StandardErrorPath</key>
  <string>{{REPO}}/server/.logs/static.err.log</string>

  <key>WorkingDirectory</key>
  <string>{{REPO}}/site</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
